# command-service/commandrunner.service
[Unit]
Description=Homepage Command Runner

[Service]
ExecStart=/usr/bin/COMMANDRUNNER.sh

[Install]
WantedBy=multi-user.target